# backend/allmon3-cleanup.service
[Unit]
Description=Allmon3 Archive Cleanup Daemon
After=network.target

[Service]
Type=simple
ExecStart=VENV_PYTHON INSTALL_DIR/cleanup_daemon.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
